# ModBus2TCPIP/src/inittab
#inet:respawn:/bin/inetd
boa:respawn:/bin/boa
